# installer/litelayer-cloudflare.service
[Unit]
Description=LiteLayer Cloudflare Tunnel
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/
After=network-online.target litelayer.service
Wants=network-online.target
Requires=litelayer.service

[Service]
Type=simple
# Default: free quick tunnel (random *.trycloudflare.com URL). A named tunnel set
# up from the UI drops a /etc/litelayer/cloudflare.env that overrides CF_TUNNEL_ARGS
# with `run --token <token>` for a stable custom domain.
Environment=CF_TUNNEL_ARGS=--url http://localhost:8000
EnvironmentFile=-/etc/litelayer/cloudflare.env
ExecStart=/usr/bin/cloudflared tunnel $CF_TUNNEL_ARGS
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=litelayer-cloudflare

[Install]
WantedBy=multi-user.target
